# furniture-repricer-scheduler.service
[Unit]
Description=Furniture Repricer Scheduler
After=network.target

[Service]
Type=simple
User=repricer
Group=repricer
WorkingDirectory=/home/repricer/furniture-repricer
Environment="PATH=/home/repricer/furniture-repricer/venv/bin:/usr/bin:/bin"
ExecStart=/home/repricer/furniture-repricer/venv/bin/python run_scheduler.py
Restart=always
RestartSec=30
StandardOutput=append:/home/repricer/furniture-repricer/logs/scheduler.log
StandardError=append:/home/repricer/furniture-repricer/logs/scheduler-error.log

# Resource limits (optional but recommended)
MemoryLimit=6G
CPUQuota=200%

# Restart limits
StartLimitInterval=10min
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
